# utils/spotify-bot.service
[Unit]
Description=Spotify Threads Bot
After=network.target

[Service]
Type=simple
User=votre_utilisateur
WorkingDirectory=/home/votre_utilisateur/spotify
# Utilisation directe de l'interpréteur Python du venv
ExecStart=/home/votre_utilisateur/spotify/venv/bin/python run_bot.py
Restart=always
RestartSec=10

# Logs
StandardOutput=append:/home/votre_utilisateur/spotify/logs/service.log
StandardError=append:/home/votre_utilisateur/spotify/logs/service.err

# Limites ressources
MemoryLimit=500M

[Install]
WantedBy=multi-user.target
